# config/crontab
0 8 * * * php /var/www/html/app/api/Cron/check_notification.php
